# systemd/wallpaper.service
[Unit]
Description=Choose a random wallpaper to use with AWWW

[Service]
ExecStart=/home/mabel/.scripts/wallpaper.sh
